# xfce4-dotfiles/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Redshift
Comment=Night Light
Path=/usr/bin/redshift
Exec=redshift -O 2800
Terminal=True